# jupyter.conf
[supervisorctl]

[program:jupyter]
command=/home/bob/anaconda3/bin/jupyter-lab
redirect_stderr=true
directory=/home/bob
user=bob
autostart=true
autorestart=true
stdout_logfile=/tmp/jupyter.log

